# deploy/emperor.service
[Unit]
Description=Emperor media-server client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DEPLOY_USER__
Group=__DEPLOY_USER__
WorkingDirectory=__DEPLOY_PATH__
ExecStart=__NODE_BIN__ build
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=HOST=__APP_HOST__
Environment=PORT=__APP_PORT__
Environment=ORIGIN=__APP_ORIGIN__

# Harden lightly; deploy user owns the app tree.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
