# deploy/emperor.service
[Unit]
Description=Emperor media-server client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__DEPLOY_PATH__
ExecStart=__NODE_BIN__ build
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=HOST=__APP_HOST__
Environment=PORT=__APP_PORT__
Environment=ORIGIN=__APP_ORIGIN__

[Install]
WantedBy=default.target
